# pi/services/mosque-kiosk.service
[Unit]
Description=Minara Kiosk (Pi 4 Wayland)
After=wifi-setup.service

[Service]
User=kiosk
Group=kiosk
PAMName=login
TTYPath=/dev/tty1
StandardInput=tty
StandardOutput=journal
StandardError=journal

Environment=XDG_RUNTIME_DIR=/run/user/1001
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1001/bus

ExecStartPre=/bin/sleep 3
ExecStart=/opt/mosque/scripts/start-kiosk.sh

Restart=always
RestartSec=5

# start-kiosk.sh zählt Abstürze selbst und rebootet bei CRASH_LIMIT=10
# kein StartLimitBurst hier — der würde den Script-internen Reboot verhindern
StartLimitIntervalSec=0

TimeoutStartSec=60
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
