# pi/services/mosque-wifi-watchdog.service
[Unit]
Description=Mosque WiFi Watchdog
After=wifi-setup.service

[Service]
Type=simple
ExecStart=/opt/mosque/scripts/wifi-watchdog.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
